# files/usr/etc/ld.so.preload
/usr/lib64/glibc-hwcaps/x86-64-v4/libhardened_malloc.so
